# templates/ld_auto_patching/ld-apply-patches.service.erb
[Unit]
Description=Automatically apply LD patches
After=local-fs.target

[Service]
Type=oneshot
ExecStart=<%= @localbindir %>/ld-apply-patches
KillMode=process
TimeoutStopSec=900
