# templates/hello_world_flask/hello-world-flask.service.erb
[Unit]
Description=Hello World application in Flask
Wants=basic.target
After=basic.target network.target

[Service]
<% if @facts['os']['family'] == 'RedHat' -%>
ExecStart=gunicorn --bind <%= @ip_address %>:<%= @port %> wsgi:app
<% else -%>
ExecStart=python3 -m gunicorn --bind <%= @ip_address %>:<%= @port %> wsgi:app
<% end -%>
KillMode=process
WorkingDirectory=<%= @app_dir %>

[Install]
WantedBy=multi-user.target

